# .config/systemd/user/snixembed.service
[Unit]
Description=snixembed

[Service]
ExecStart=/usr/bin/snixembed

[Install]
WantedBy=default.target
