# .config/autostart/Proton Mail Bridge.desktop
[Desktop Entry]
Type=Application
Name=Proton Mail Bridge
Exec="/usr/bin/protonmail-bridge" "--no-window"
X-GNOME-Autostart-enabled=true
